# tools/team-chat/team-chat-heartbeat@.service
[Unit]
Description=Team Chat Heartbeat for %i
After=network.target

[Service]
Type=simple
User=aiciv
Environment=HOME=/home/aiciv
Environment=ROOM_CONFIG=/home/aiciv/duo/room-config.json
Environment=DUO_DIR=/home/aiciv/duo
WorkingDirectory=/home/aiciv
ExecStart=/usr/bin/python3 /home/aiciv/purebrain_portal/tools/team-chat/room_heartbeat.py
Restart=always
RestartSec=10
StandardOutput=append:/home/aiciv/duo/logs/room_heartbeat.log
StandardError=append:/home/aiciv/duo/logs/room_heartbeat.log

[Install]
WantedBy=multi-user.target
