# tools/team-chat/team-chat-poller@.service
[Unit]
Description=Team Chat Message Poller for %i
After=network.target
Wants=team-chat-heartbeat@%i.service

[Service]
Type=simple
User=aiciv
Environment=HOME=/home/aiciv
Environment=ROOM_CONFIG=/home/aiciv/duo/room-config.json
Environment=DUO_DIR=/home/aiciv/duo
WorkingDirectory=/home/aiciv
ExecStart=/usr/bin/python3 /home/aiciv/purebrain_portal/tools/team-chat/room_poller.py
Restart=always
RestartSec=10
StandardOutput=append:/home/aiciv/duo/logs/room_poller.log
StandardError=append:/home/aiciv/duo/logs/room_poller.log

[Install]
WantedBy=multi-user.target
